# .config/systemd/user/system-update-refresh.timer
[Unit]
Description=Periodically refresh system update cache

[Timer]
OnBootSec=2min
OnUnitActiveSec=30min
RandomizedDelaySec=5min

[Install]
WantedBy=timers.target
